# crontab
0 9 * * * python3 ~/git/ipchecker/checker.py
